# roles/minecraft-world-backup/templates/minecraft-backup.cron.j2
#10 */2 * * *  backup  /home/backup/minecraft-world-backup.sh
10 05-17/4,19-23/2,01-03/2 * * * backup /home/backup/minecraft-backup-logging.sh
