# roles/logger/files/domain-notifier.cron
*/5 * * * *  logger  /home/logger/domain-notifier.sh
