# roles/logger/files/temp-logger.cron
*/3 * * * *  logger  /home/logger/temp-logger.sh
