# roles/monitoring/files/cpu-mon.cron
*/3 * * * *  monitoring  /home/monitoring/cpu-mon.sh
