# roles/monitoring/files/ngrok-mon.cron
*/5 * * * *  monitoring  /home/monitoring/ngrok-mon.sh
