# roles/monitoring/files/temp-mon.cron
*/3 * * * *  monitoring  /home/monitoring/temp-mon.sh
